# systemd/subconverter.service
[Unit]
Description=Subconverter Service

[Service]
ExecStart=subconverter
Restart=always

[Install]
WantedBy=multi-user.target
